# service/launchd.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>{{LABEL}}</string>
  <key>ProgramArguments</key>
  <array>
    <string>{{NODE}}</string>
    <string>{{SERVER}}</string>
    <string>run</string>
    <string>--no-self-update</string>
    <string>--profile</string>
    <string>{{PROFILE}}</string>
    <string>--port</string>
    <string>{{PORT}}</string>
    <string>--host</string>
    <string>{{HOST}}</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>{{HOME}}</string>
    <key>CLAUDE_AGENT_API_NO_SELF_UPDATE</key>
    <string>1</string>
    <key>PATH</key>
    <string>{{PATH}}</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>{{STDOUT}}</string>
  <key>StandardErrorPath</key>
  <string>{{STDERR}}</string>
</dict>
</plist>
