# service/systemd.service.template
[Unit]
Description=claude-agent-api bridge (%i)
After=network-online.target

[Service]
Type=simple
ExecStart={{NODE}} {{SERVER}} run --no-self-update --profile {{PROFILE}} --port {{PORT}} --host {{HOST}}
Environment=CLAUDE_AGENT_API_NO_SELF_UPDATE=1
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
